# scripts/com.aibox.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aibox.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/daniel/GitHub/AI-Box/.venv/bin/python</string>
        <string>/Users/daniel/GitHub/AI-Box/scripts/backup_all.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/daniel/GitHub/AI-Box</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/daniel/GitHub/AI-Box/logs/backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/daniel/GitHub/AI-Box/logs/backup.error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
